# confs/dotprefix/config/systemd/user/clash.target
[Unit]
Description=Clash, a tool fucking the GFW
Requires=systemd-networkd.server
